# earthnow.slice
[Unit]
Description=Limit resources for EarthNow Background
DefaultDependencies=no
Before=slices.target
[Slice]
CPUQuota=10%
MemoryLimit=1G
